# deploy/data-cleaner.service
[Unit]
Description=Data Cleaner FastAPI Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/python-homework
ExecStart=/root/python-homework/.venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
